# installation/ansible/cluster/30.spark/plays/roles/common/templates/profile.d/spark.sh.j2
export SPARK_HOME="{{ SPARK_HOME }}"
export SPARK_CONF_DIR="${SPARK_HOME}/conf"
export PATH="$PATH:${SPARK_HOME}/bin:${SPARK_HOME}/sbin"